# core/observer-api/quantlab-observer.service
[Unit]
Description=QuantLab Observer API
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/quantlab/api/observer-api
Environment=OBSERVER_MODE=1
Environment=OBSERVER_PORT=3000
# Using absolute path to node and tsx
ExecStart=/usr/bin/node /home/deploy/quantlab/api/observer-api/node_modules/.bin/tsx index.ts
Restart=always
RestartSec=5
StandardOutput=append:/home/deploy/quantlab/api/observer-api/server.log
StandardError=append:/home/deploy/quantlab/api/observer-api/server.log

[Install]
WantedBy=multi-user.target
