# core/compressor/quantlab-compact.service
[Unit]
Description=QuantLab Daily Parquet Compaction Job
After=network.target

[Service]
Type=oneshot
User=deploy
WorkingDirectory=/home/deploy/quantlab/api/core
ExecStart=/usr/bin/python3 /home/deploy/quantlab/api/core/compressor/run.py
StandardOutput=journal
StandardError=journal
Restart=no
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
